# crontabs/root
#!/bin/sh
*/10 * * * * 	/bin/sh /radoneye-mqtt-bridge/radoneye-mqtt-bridge.sh
